# sentinel.service
[Unit]
Description=Sentinel - Advanced Abuse Detection by Novel
Documentation=https://github.com/your-repo/sentinel
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/etc/sentinel
ExecStart=/usr/local/bin/sentinel --daemon --log-level info
ExecStop=/usr/local/bin/sentinel --action stop
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentinel
TimeoutStopSec=10
KillMode=mixed

[Install]
WantedBy=multi-user.target